# system/tt-dialer3307.service
[Unit]
Description=ticket assign manager service
After=network.target

[Service]
User=root
WorkingDirectory=/Czentrix/apps/TicketMgr
ExecStart=/Czentrix/apps/cc_app_api/venv/bin/python ticketAssignMngr.py 3307 '/var/lib/mysql3307/mysql3307.sock'

[Install]
WantedBy=multi-user.target
